# downloads/code/cron-replacement/timer-hourly.target
[Unit]
Description=Hourly Timer Target
StopWhenUnneeded=yes
