# static/downloads/code/cron-replacement/pkgstats.service
[Unit]
Description=Run pkgstats

[Service]
User=nobody
ExecStart=/usr/bin/pkgstats
